# supervisord.conf
[supervisord]
nodaemon=true

[program:boot]
command=/bin/bash -c "cd /app; boot prod"
